# deploy/supervisor/subscription-service.conf
[program:subscription-service]
directory=/Users/gaoyong/Documents/work/xinyuan_tech/subscription-service
command=/Users/gaoyong/Documents/work/xinyuan_tech/subscription-service/bin/server -conf /Users/gaoyong/Documents/work/xinyuan_tech/subscription-service/configs/config.yaml
autostart=true
autorestart=true
startsecs=10
startretries=3
user=gaoyong
redirect_stderr=true
stdout_logfile=/Users/gaoyong/Documents/work/xinyuan_tech/subscription-service/logs/subscription-service.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=GIN_MODE=release
